# antcam.service
[Unit]
Description=Count the number of ants in view
[Service]
Type=oneshot
ExecStart=/app/antcam
